# auroville-connect.service
[Unit]
Description=Auroville Connect Platform
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/AurovilleConnect
Environment=NODE_ENV=production
Environment=PORT=5000
Environment=HOST=0.0.0.0
Environment=DATABASE_URL=postgresql://postgres:postgres@localhost:5432/auroville
Environment=VITE_API_URL=https://api.auroville.social
Environment=VITE_APP_URL=https://auroville.social
Environment=FRONTEND_URL=https://auroville.social
Environment=API_URL=https://api.auroville.social

# Install dependencies for both frontend and server
ExecStartPre=/usr/bin/npm install
ExecStartPre=/bin/bash -c 'cd server && /usr/bin/npm install'

# Build frontend
ExecStartPre=/usr/bin/npm run build:prod

# Start the server
WorkingDirectory=/root/AurovilleConnect/server
ExecStart=/usr/bin/npm start

# Restart on failure
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
